# package/deb/netmicroscope-daemon/DEBIAN/postinst
#!/bin/bash

chmod +x /etc/init.d/netmicroscope-daemon
update-rc.d netmicroscope-daemon defaults
